# etc/udev/rules.d/99-input.rules
KERNEL=="input[0-9]*", RUN+="/root/bin/a2dp-autoconnect"
